# infra/systemd/lotsman-watchdog.timer
[Unit]
Description=Run Лоцман pod watchdog every 2 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=2min
AccuracySec=15s

[Install]
WantedBy=timers.target
